# toyserver.service.in
[Unit]
Description=Toy Racer server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=INSTALL_USER
WorkingDirectory=LOCALSTATEDIR/lib/toyserver
ExecStart=SBINDIR/toyserver
StandardOutput=append:LOCALSTATEDIR/log/toyserver.log
LimitCORE=200M

[Install]
WantedBy=multi-user.target
